# launchd/com.mutaafaziz.market-digest-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mutaafaziz.market-digest-weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mutaafaziz/Desktop/projects/market-digest/.venv/bin/python</string>
        <string>/Users/mutaafaziz/Desktop/projects/market-digest/scripts/run_digest.py</string>
        <string>--type</string>
        <string>weekly</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mutaafaziz/Desktop/projects/market-digest</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/mutaafaziz/Desktop/projects/market-digest</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <!-- Friday at 5:30 PM -->
        <key>Weekday</key><integer>5</integer>
        <key>Hour</key><integer>17</integer>
        <key>Minute</key><integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/mutaafaziz/Desktop/projects/market-digest/logs/weekly_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mutaafaziz/Desktop/projects/market-digest/logs/weekly_stderr.log</string>
</dict>
</plist>
